# cron.d/root
0 0 * * 6 /usr/local/bin/auto_backup
55 23 * * * /usr/local/bin/culvert_score_update_reminder
